# files/etc/xdg/autostart/gsm.desktop
[Desktop Entry]
Name=GSM
Comment=GSM
Exec=gnome-settings-daemon
Terminal=false
Type=Application
StartupNotify=false
OnlyShowIn=AWESOME;
