# files/etc/xdg/autostart/sound.desktop
[Desktop Entry]
Name=GSA
Comment=GSA
Exec=gnome-sound-applet
Terminal=false
Type=Application
StartupNotify=false
OnlyShowIn=AWESOME;
